# release/config/inhive-core.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="inhive-core service"
supervisor="supervise-daemon"
command="/usr/bin/inhive-core"
command_args="-D /var/lib/inhive-core -C /etc/inhive-core run"
extra_started_commands="reload"

depend() {
	after net dns
}

reload() {
	ebegin "Reloading $RC_SVCNAME"
        $supervisor "$RC_SVCNAME" --signal HUP
	eend $?
}